# deploy/babylon-ingest-wsl2.service
[Unit]
Description=Babylon Ingest Daemon (Phase 2 — WSL2)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/mnt/b/Babylon/app/ingest
EnvironmentFile=/mnt/b/Babylon/app/.env
ExecStart=/mnt/b/Babylon/app/ingest/venv/bin/python3 daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=babylon-ingest
MemoryMax=4G

[Install]
WantedBy=multi-user.target
